# listener.service
[Unit]
Description=My Python Script

[Service]
Type=simple
ExecStart=/usr/bin/python3 /path/to/your_script.py
WorkingDirectory=/path/to/script/directory
Restart=always
User=your_username

[Install]
WantedBy=multi-user.target
